# meta-custom/recipes-platform-fix/rcar-pvr-watchdog/files/pvr-watchdog.service
[Unit]
Description=PVR Watchdog

[Service]
ExecStart=/usr/bin/pvr-watchdog
Type=oneshot
