# meta-agl-early/recipes-support/rcar-pvr-watchdog/files/pvr-watchdog.timer
[Unit]
Description=Run pvr-watchdog every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=basic.target

